# tunegpu.timer
[Unit]
Description=Tune GPU Settings

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target
